# .glitch/supervisord.conf
;
; For hosting on Glitch.com
;

[supervisord]
loglevel=debug

[program:bot]
command=../main.py
autorestart=true

[program:server]
command=./server.py
autorestart=true
